# www-apps/mediagoblin/files/init.d/mediagoblin
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

command="${MEDIAGOBLIN_CMD}"
command_args="${MEDIAGOBLIN_ARGS}"
pidfile="${MEDIAGOBLIN_PID}"

depend() {
    need net logger
    need postgresql
    need celery-worker

    provide mediagoblin
}

start() {

    # Make sure the log and pid file exist and are writeable
    for path in "${MEDIAGOBLIN_LOG}" "${MEDIAGOBLIN_PID}" ; do
        dir=$(dirname $path)
        sudo -u ${MEDIAGOBLIN_USER} test ! -d "${dir}" && \
            mkdir $dir && \
            chown ${MEDIAGOBLIN_USER} "${dir}"
    done

    ebegin "Starting ${SVCNAME}"

    start-stop-daemon --start \
        --chdir ${MEDIAGOBLIN_DIR} \
        --name ${SVCNAME} \
        --exec ${command} \
        --user ${MEDIAGOBLIN_USER}:${MEDIAGOBLIN_GROUP} \
        --pidfile ${pidfile} \
        -- ${command_args}
    
    eend $?
}
