# net-analyzer/prometheus/files/prometheus-initd
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

exec=/usr/bin/prometheus
config=/etc/prometheus.yml

depend() {
	need net
}

start() {
	ebegin "Starting Prometheus Server"
	start-stop-daemon --start --quiet --exec ${exec} --\
		-config.file=${config}
	eend $?
}

stop() {
	ebegin "Starting Prometheus Server"
	start-stop-daemon --stop --quiet --exec ${exec}
	eend $?
}
