# net-analyzer/prometheus-glsa-exporter/files/prometheus-glsa-exporter-initd
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Prometheus Node Exporter"

pidfile="/run/${SVCNAME}.pid"
command=/usr/bin/prometheus-glsa-exporter
command_args="${PROMETHEUS_GLSA_EXPORTER_ARGS}"
command_background="true"

start_stop_daemon_args="-u ${PROMETHEUS_GLSA_EXPORTER_USER} -l ${PROMETHEUS_GLSA_EXPORTER_LISTEN}"

depend() {
	need net
}
